# rootfs/etc/aura/services.d/health-monitor.service
[Service]
Name=health-monitor
Description=Service health monitor and self-repair
Module=services.health_monitor:HealthMonitor
AutoStart=true
RestartPolicy=on-failure
RestartDelay=10
DependsOn=job-queue,logging

[Health]
CheckInterval=30
MaxFailures=5
