# rootfs/etc/aura/services.d/module-builder.service
[Service]
Name=module-builder
Description=Module scaffolding and self-expansion engine
Module=services.module_builder:ModuleBuilder
AutoStart=false
RestartPolicy=never
DependsOn=logging

[Health]
CheckInterval=60
MaxFailures=1
