# rootfs.pre_rollback/etc/aura/services.d/network.service
[Service]
Name=network
Description=Network connectivity monitor
Module=services.network_service:NetworkService
AutoStart=true
RestartPolicy=on-failure
RestartDelay=5
DependsOn=

[Health]
CheckInterval=30
MaxFailures=3
